# services/clamav-freshclam/run
#!/bin/sh
sleep 6h
echo "Checking for ClamAV database updates"
/usr/bin/freshclam
