# services/memcached/run
#!/bin/bash
exec memcached -u memcache -p 11211 -l 127.0.0.1  
