# services/clamav/run
#!/bin/sh

if [ ! -d /var/run/clamav ]; then
  mkdir -p /var/run/clamav
  chown clamav:root /var/run/clamav
  chmod 755 /var/run/clamav
fi

exec /usr/sbin/clamd -F
