# services/clamav/run
#!/bin/sh

# Refresh ClamAV cvd files
echo "Starting ClamAV..."
exec /usr/sbin/clamd -F
