# services/clamav/run
#!/bin/sh

# Refresh ClamAV cvd files
echo "Starting ClamAV..."
exec /usr/sbin/clamd -c /etc/clamd.d/amavisd.conf -F
